# rtserver/tmp/rc.local
#!/bin/sh

su - ubuntu -c "/home/ubuntu/bin/start_rtserver.sh"

su - ubuntu -c "/home/ubuntu/bin/update-public-dns.sh rtserver"
su - ubuntu -c "/home/ubuntu/bin/update-local-dns.sh rtserver"

exit 0
